# packaging/com.activitywatch.aw-watcher-cmux.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd service for aw-watcher-cmux (spec §12).

  ⚠️  KNOWN LIMITATION: cmux's control socket uses access_mode 'cmuxOnly' and
  only authorizes callers running INSIDE a cmux surface. A launchd LaunchAgent
  is detached and has no caller surface, so cmux rejects every query (SIGPIPE)
  and the watcher records nothing. Verified empirically. This plist is kept for
  reference / future use once external socket auth is configured in cmux. For
  now, run aw-watcher-cmux from inside a cmux tab instead (see README).

  Install:
    1. Edit the path in <string> below to point at your installed entry point
       (e.g. ~/.local/bin/aw-watcher-cmux from `pipx install .`).
    2. cp this file to ~/Library/LaunchAgents/
    3. launchctl load ~/Library/LaunchAgents/com.activitywatch.aw-watcher-cmux.plist

  Prefer aw-qt management if you run it: aw-qt auto-discovers executables named
  aw-watcher-* and can start/stop this watcher alongside the others.
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.activitywatch.aw-watcher-cmux</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/CHANGE_ME/.local/bin/aw-watcher-cmux</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!--
      Restart only on a crash (non-zero exit), not on a clean exit. With a bare
      <true/> a misconfigured path or startup crash would respawn in a tight
      loop; SuccessfulExit=false lets a clean shutdown stay down. ThrottleInterval
      caps the respawn rate as a backstop.
    -->
    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/CHANGE_ME/Library/Logs/activitywatch/aw-watcher-cmux.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/CHANGE_ME/Library/Logs/activitywatch/aw-watcher-cmux.log</string>

    <key>ProcessType</key>
    <string>Background</string>
  </dict>
</plist>
